# deploy/linux/edex.service
[Unit]
Description=eDEX-UI LAN browser service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/edex-ui/edex.env
User=edex
Group=edex
WorkingDirectory=/opt/edex-ui
ExecStart=/usr/local/bin/run-edex-session-linux.sh
Restart=always
RestartSec=5
TimeoutStopSec=20
KillMode=control-group

[Install]
WantedBy=multi-user.target

